# rules/persistence/startup_folder_file_drop.yml
title: File Dropped in Windows Startup Folder
id: 7a3f2e01-d4c8-4b91-a1f3-5e8c9d2b4a04
status: experimental
description: |
    Detects file creation events in Windows startup folders (both per-user and all-users).
    Any executable, script, or shortcut placed in these directories will run automatically
    at user logon. This is a low-sophistication but highly effective persistence technique
    used by commodity malware, ransomware droppers, and even some APT groups as a
    quick-and-dirty foothold. Legitimate software rarely drops files here — most use
    registry Run keys or scheduled tasks instead.
references:
    - https://attack.mitre.org/techniques/T1547/001/
author: jawsec
date: 2026/04/01
modified: 2026/04/01
tags:
    - attack.persistence
    - attack.t1547.001
logsource:
    category: file_event
    product: windows
detection:
    selection:
        TargetFilename|contains:
            - '\Start Menu\Programs\Startup\'
            - '\AppData\Roaming\Microsoft\Windows\Start Menu\Programs\Startup\'
            - '\ProgramData\Microsoft\Windows\Start Menu\Programs\Startup\'
    filter_extensions:
        TargetFilename|endswith:
            - '.ini'
            - '.desktop.ini'
    condition: selection and not filter_extensions
falsepositives:
    - Legitimate applications that add themselves to startup during installation
    - IT administrators deploying logon scripts via startup folder (uncommon but possible)
    - OneDrive or cloud sync restoring previously synced startup items
level: high
# ---
# TESTING METHODOLOGY
# To trigger this rule in a lab environment:
#   1. Ensure Sysmon is installed with file creation events enabled (Event ID 11)
#   2. Open PowerShell and run:
#      New-Item -Path "$env:APPDATA\Microsoft\Windows\Start Menu\Programs\Startup\test_persistence.bat" -Value "echo test"
#   3. Verify Sysmon Event ID 11 captures the file creation
#   4. Confirm your SIEM fires the alert
#   5. Clean up: Remove-Item "$env:APPDATA\Microsoft\Windows\Start Menu\Programs\Startup\test_persistence.bat"
# Expected result: Alert fires on any file creation in the Startup directory
